# config/cron/feed.cron
# Queue feed updates.
2 * * * * cd /var/www/oursignal && ./bin/oursignal-feed update
