# config/cron/update.cron
# Required environment stuff. How does Barney do this?
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Schedule an update every 15 minutes.
*/15 * * * * cd /home/oursignal/oursignal && ./bin/oursignal-update all
