# config/cron/update.cron
# Schedule an update every 5 minutes.
*/5 * * * * cd /home/oursignal/oursignal && ./bin/oursignal-update all
