# config/cron/update.cron
# Schedule an update every 5 minutes.
*/5 * * * * cd /var/www/oursignal && ./bin/oursignal-update all
